# conf.d/polipo.conf
[program:polipo]
command=/bin/bash /start.sh polipo
autostart=true
autorestart=true
