# conf.d/shadowsocks.conf
[program:sslocal]
command=/bin/bash /start.sh sslocal
autostart=true
autorestart=true
